# modules/portaudit/files/portaudit.cron
#!/bin/sh

# Managed by Puppet

echo
echo "Checking packages for vulnerabilities:"
out=$(portaudit -Faq | tee /dev/stderr | wc -l)
[ "$out" -gt 0 ] && exit 1
exit 0
