# modules/postfix/files/heartbeat.cron
#!/bin/sh

# Managed by Puppet

date | mail -s "heartbeat `hostname -s`" root
