# modules/fsnapshot/files/daily.cron
#!/bin/sh

# Managed by Puppet

periodic-snapshot daily
