# modules/portsnap/files/portsnap.cron
#!/bin/sh

# Managed by Puppet

echo
echo "Updating ports: "

portsnap cron
[ "$n" -gt 0 ] && rc=1
if [ -e /usr/ports/.portsnap.INDEX ]; then
  portsnap update
else
  portsnap extract
fi
[ "$n" -gt 0 ] && rc=1

if [ -x /usr/local/sbin/portsdb ]; then
  portsdb -u
  [ "$n" -gt 0 ] && rc=1
fi

exit ${rc:-0}
