# old/ports/files/ports-clean.cron
#!/bin/sh

# Managed by Puppet

echo
echo 'Cleaning up ports tree:'
/usr/local/sbin/portsclean -C -DD > /dev/null || exit 1
