# modules/package/files/pkgaudit.cron
#!/bin/sh

# Managed by Puppet

echo
echo "Checking packages for vulnerabilities:"
out=$(pkg audit -Fq | tee /dev/stderr | wc -l)
[ "$out" -gt 0 ] && exit 1
exit 0
