# modules/fupdate/files/freebsd-update.cron
#!/bin/sh

# Managed by Puppet

/usr/sbin/freebsd-update cron
