# modules/backup/files/client/nodumps.cron
#!/bin/sh

# Managed by Puppet

echo
echo "Setting nodump flags:"
/usr/local/sbin/nodumps
