# modules/zfs/files/zfs-backup.cron
#!/bin/sh

# Managed by Puppet

set -e

echo
echo "Backing up ZFS properties:"

cd /var/backups
if [ -e zfs.props ]; then
  mv zfs.props zfs.props.old
fi
zfs list -o all > zfs.props
