# modules/syslog/files/newsyslog.cron
#!/bin/sh

# Managed by Puppet

echo
echo "Creating and/or trimming log files:"
newsyslog -C
