# old/backup/files/client/backup.cron
#!/bin/sh

# Managed by Puppet

echo
echo "Starting backup:"
/usr/local/sbin/backup
